# ansible/templates/market.service.j2
[Unit]
Description=Plebeian Market (Bun Fullstack)
After=network.target nak.service
Requires=nak.service

[Service]
Type=simple
User={{ app_user }}
WorkingDirectory={{ app_dir }}
Environment=PATH=/home/{{ app_user }}/.bun/bin:/usr/bin:/bin
EnvironmentFile=-{{ app_dir }}/.env
ExecStart=/bin/bash -c "cd {{ app_dir }} && /home/{{ app_user }}/.bun/bin/bun run startup || true && /home/{{ app_user }}/.bun/bin/bun start"
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
